# jarvis-dashboard-static.service
[Unit]
Description=Jarvis War Room Dashboard Frontend (SPA static server)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/.hermes/profiles/jarvis/plugins/jarvis-dashboard
Environment="PATH=/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/ubuntu/.hermes/profiles/jarvis/plugins/jarvis-dashboard/spa_server.py 8503 /home/ubuntu/.hermes/profiles/jarvis/plugins/jarvis-dashboard/frontend/public
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
